# files/system/usr/lib/udev/rules.d/99-deepcool-digital.rules
# DeepCool Digital Linux - udev rules for rootless operation
# Intel RAPL energy monitoring interface
SUBSYSTEM=="powercap", KERNEL=="intel-rapl:*", RUN+="/usr/bin/chmod -R 0644 /sys/devices/virtual/powercap/intel-rapl/"

# DeepCool USB HID devices (adjust idVendor and idProduct based on your specific hardware)
# Common DeepCool vendor ID: 3633 (0x0e31)
SUBSYSTEM=="usb", ATTRS{idVendor}=="3633", MODE="0666"
SUBSYSTEM=="hidraw", ATTRS{idVendor}=="3633", MODE="0666"
